# spr/cron/expire_counters.cron
*/10 * * * * python /app/scripts/expire_counters.py >> /var/log/cron.log 2>&1
